# php-fpm.ini
[program:php-fpm]
command=/usr/bin/php-fpm5 -F
