# php-fpm.ini
[program:php-fpm]
command=/usr/bin/php-fpm -F
